# ip_update_notifier/check_ip.service
[Unit]
Description=Check public IP and send email if changed

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_ip.sh
